# scripts/prod/expenseapp.service
[Unit]
Description=ExpenseApp
After=network.target

[Service]
Type=simple
User=expenseapp
Group=expenseapp
WorkingDirectory=/opt/expenseapp/
ExecStart=/opt/expenseapp/run_as_service.sh

# Security best practices (Optional but recommended)
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target